# src/books.conf
[program:books]
command=
directory=/home/ubuntu/myproject
autostart=true
autorestart=true
stderr_logfile=/var/log/myproject.err.log
stdout_logfile=/var/log/myproject.out.log